# bin/templates/gunicorn.service.template
# confirm the working directory

[Unit]
Description=gunicorn {site} daemon
Requires=gunicorn.{site}.socket
After=network.target

[Service]
User=edc
Group=www-data
WorkingDirectory=/home/edc/app
ExecStart=/home/edc/miniconda3/envs/edc/bin/gunicorn \
    --access-logfile - \
    --workers 3 \
    --timeout=300 \
    --bind unix:/run/gunicorn.{site}.sock \
    inte_edc.wsgi.{site}:application

[Install]
WantedBy=multi-user.target