# bin/systemd/gunicorn.socket.template
[Unit]
Description=gunicorn {site} socket

[Socket]
ListenStream=/run/gunicorn-{site}.sock

[Install]
WantedBy=sockets.target